# genome-www/etc/cron.daily/genome-task-cleanup.sh
#!/bin/sh

find /var/cache/genome/task_runner_output/ -ctime +7 -exec rm -f {} \;
